# dev-embedded/libzinstr/files/zidataserver
#!/sbin/openrc-run

command=/usr/bin/ziDataServer
MYPROG_OPTS=""
description="Zurich Instruments Data Server"
pidfile=/var/run/zi_dataserver.pid


depend() {
   need net
}

start() {
   ebegin "Starting ${description}"
   start-stop-daemon --start --quiet --background --make-pidfile --pidfile ${pidfile} \
     --exec ${command} \
     -- ${MYPROG_OPTS}
   eend $?
}

stop() {
   ebegin "Stopping ${description}"
   start-stop-daemon --stop --quiet --pidfile ${pidfile}
   eend $?
}

